# STI4A/linux_admin/script-init-exemple
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

opts="${opts}"

depend() {
	need net
	after mysql
}

start() {
	ebegin "Starting  Open Source Security Information Management: OSSIM"
		start-stop-daemon --start --exec /usr/bin/ossim-agent -- -f -d 
		eend $?
}

stop() {
	ebegin "Stoping Open Source Security Information Management: OSSIM"
	kill -9 /var/lib/ossim/ossim.pid
	eend $?

}

